# ops/launchd/com.hakanssonlabs.newsfeed.pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  User LaunchAgent template for the Newsfeed producer pipeline.
  Render the __PLACEHOLDERS__ at install time; see ops/launchd/README.md.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hakanssonlabs.newsfeed.pipeline</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO_DIR__/ops/launchd/run-pipeline.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>21600</integer>

  <key>WorkingDirectory</key>
  <string>__REPO_DIR__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <!-- Include Node, npm/npx, Hermes, and git for launchd's minimal environment. -->
    <key>PATH</key>
    <string>__HOME__/.hermes/hermes-agent/venv/bin:__NODE_BIN_DIR__:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>__HOME__</string>
    <key>HERMES_BIN</key>
    <string>__HOME__/.hermes/hermes-agent/venv/bin/hermes</string>
  </dict>

  <key>StandardOutPath</key>
  <string>__REPO_DIR__/ops/launchd/logs/pipeline.out.log</string>
  <key>StandardErrorPath</key>
  <string>__REPO_DIR__/ops/launchd/logs/pipeline.err.log</string>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
